# infra/ansible/roles/sls_bot/templates/sls-backup.timer.j2
[Unit]
Description=Backups Restic SLS Bot

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
