# scripts/deploy/systemd/sls-bot.service
[Unit]
Description=SLS Trading Bot (Bybit bridge)
After=network.target

[Service]
Type=simple
User={{SVC_USER}}
Group={{SVC_USER}}
WorkingDirectory={{APP_ROOT}}/bot
EnvironmentFile=-/etc/sls_bot.env
Environment="PYTHONPATH={{APP_ROOT}}/bot"
Environment="SLSBOT_CONFIG=${SLSBOT_CONFIG:-{{APP_ROOT}}/config/config.json}"
ExecStart={{APP_ROOT}}/venv/bin/uvicorn sls_bot.app:app --host 0.0.0.0 --port ${SLS_BOT_PORT:-8080}
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
